# sv/sonar/run
#!/bin/sh

sv start mysql || exit 1

exec 2>&1
exec /sonar/bin/linux-x86-64/sonar.sh console
